# supervisord.conf
[supervisord]
user=root
nodaemon=false
logfile=/etc/supervisord.log
childlogdir=/etc
pidfile = /etc/supervisord.pid

[program:resapi-smoke-test]
command=/resapi-smoke-test
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=0
